# common/usr/lib/systemd/user/swayidle-restart.path
[Unit]
Description=Automatically restart the swayidle when its configuration changes
PartOf=swayidle.service

[Path]
PathChanged=%h/.config/swayidle/config

[Install]
WantedBy=paths.target
