# common/usr/lib/systemd/user/swayidle-restart.service
[Service]
ExecStart=systemctl --user restart swayidle.service
Type=oneshot

[Install]
WantedBy=graphical-session.target
