# common/usr/lib/systemd/user/waybar-restart.service
[Service]
ExecStart=systemctl --user restart waybar.service
Type=oneshot
