# data/ubuntusway/skel/.config/autostart/ubuntu-sway-welcome.desktop
[Desktop Entry]
Name=Ubuntu Sway Welcome
Exec=ubuntu-sway-welcome
Icon=ubuntu-sway.png
Type=Application
Terminal=false
Categories=Settings;
