# common/etc/skel/.config/autostart/xdg-user-dirs.desktop
[Desktop Entry]
Exec=xdg-user-dirs-update
NoDisplay=true
StartupNotify=false
Terminal=false
TryExec=xdg-user-dirs-update
Type=Application
